# restic_backup.service
[Unit]
Description=Backup service file to run bash file

[Service]
ExecStart=$SH_PATH

[Install]
WantedBy=multi-user.target